# src/install/booteject.timer
[Unit]
Description=BootEject Timer

[Timer]
OnStartupSec=10

[Install]
WantedBy=timers.target
