# src/install/booteject.service
[Unit]
Description=BootEject

[Service]
ExecStart=/opt/booteject/booteject.sh
Type=oneshot
RemainAfterExit=true
